# roles/ginlong/templates/ginlong.scraper.service
[Unit]
Description=My ginlong.scraper Service
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/bin/python /home/pi/ginlong.scraper.py {{ username }} {{ password }}

[Install]
WantedBy=multi-user.target